# deployment/systemd/boggler.trigger.path
[Unit]
Description=Watch for Boggler instance

[Path]
PathChanged=<project repo path>/src/
PathChanged=<project repo path>/composer.lock
PathChanged=<project repo path>/php_cache/

[Install]
WantedBy=multi-user.target
